# linux/.zshrc.d-bak/config.d/cron/cron.daily/clamscan
#!/bin/zsh
0 20 * * * clamscan -r --quiet --bell -i /home/$USER/Downloads

