# linux/.zshrc.d-bak/config.d/cron/cron.daily/freshclam_update
#!/bin/zsh
0 19 * * * /usr/bin/freshclam --quiet
