# linux/dotfiles/.zshrc.d/config.d/cron/cron.weekly/clamscan
#!/bin/zsh
0 21 * * * clamscan -r --quiet --bell -i /home/$USER
