# linux/dotfiles/.zshrc.d/config.d/cron/cron.daily/clamscan.daily
#!/bin/zsh
clamscan -r --quiet --bell -i /home/$USER

